# home_work_25/flask-app.service
[Unit]
Description=Flask app service
After=network.target

[Service]
WorkingDirectory=/home/artilerist47/app
ExecStart=/home/artilerist47/app/venv/bin/python run.py

[Install]
WantedBy=multi-user.target